# $HOME/.local/profile.d/aliases.sh
# general-purpose aliases (special purpose aliases should get their own file)
alias ls='ls --color=auto'
alias ll='ls -lahF'
alias grep='grep --color=auto'
alias egrep='egrep --color=auto'
alias fgrep='fgrep --color=auto'
alias dir='dir --color=auto'
alias vdir='vdir --color=auto'

