# $HOME/.local/profile.d/defaults.sh
# variables to set default behavior for various apps

# editors
export VISUAL=nano
export EDITOR="$VISUAL"

# sensible-utils
export BROWSER=firefox

