# $HOME/.local/profile.d/docker.sh
# docker-related variables and aliases
export DOCKER_HOST=unix:///run/user/1000/docker.sock
