# $HOME/.local/etc/profile.d/aliases.sh
# general-purpose aliases (special purpose aliases should get their own file)
alias ll='ls -lahF'

